# templates/django/config/redis/redis.service.erb
[Unit]
Description=Redis service
Documentation=

# web
# /etc/systemd/system/redis.service

[Service]
Type=simple
# Environment=A=123
ExecStart=/usr/local/bin/docker-compose -f /opt/<%= @project_name %>/config/redis/redis-docker-compose.yaml up
Restart=always
RestartSec=3
# User=redis
# Group=redis

ExecStop=/bin/kill -s TERM $MAINPID

WorkingDirectory=/opt/<%= @project_name %>
ReadWriteDirectories=/opt/<%= @project_name %>

[Install]
WantedBy=multi-user.target
# Alias=redis.service
